# keystone/supervisord/uwsgi.conf
[program:uwsgi]
command=/usr/bin/uwsgi --master --plugin python --emperor /etc/uwsgi
stdout_events_enabled=true
stderr_events_enabled=true
